# supervisord.conf
[supervisord]
nodaemon=true

[program:jupyter notebook]
command=jupyter notebook
